# .docker/php/cron/crontab
0 12 * * * php /app/bin/console observers:notify
